# service/notifyd/run
#!/bin/sh
exec 2>&1

USER=notifyd

exec softlimit -m 134217728 \
     setuidgid $USER \
/usr/bin/notifyd